# snwbot.conf
setuid snwbot
setgid snwbot

start on startup

pre-start script
 cd /home/snwbot/bot
 npm i --production
end script

script
  cd /home/snwbot/bot
  node ./snw-bot.js
end script

respawn