# snwbot.conf
description "Snw-bot"
author  "Snekw"

setuid snwbot
setgid snwbot

start on startup
stop on shutdown

script
	exec sudo node /home/snwbot/bot/snw-bot.js
end script

respawn
respawn limit 10 5